# ai-dj.service
# systemd service file for AI DJ Docker container
# Copy to /etc/systemd/system/ai-dj.service
# Then: systemctl daemon-reload && systemctl enable ai-dj && systemctl start ai-dj

[Unit]
Description=AI DJ - Self-hosted radio station with AI commentary
Documentation=https://github.com/yourusername/ai-dj
After=docker.service
Requires=docker.service
PartOf=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/ai-dj
ExecStart=/usr/bin/docker-compose up
ExecStop=/usr/bin/docker-compose down
Restart=unless-stopped
RestartSec=10
TimeoutStopSec=60

# Run as non-root user (optional, but recommended)
# User=ai-dj
# Group=docker

# Environment
Environment=DOCKER_COMPOSE_FILE=/opt/ai-dj/docker-compose.yml

[Install]
WantedBy=multi-user.target
